# sovushka_launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.les.sovushka</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ovc/Projects/LES_v2/.venv/bin/python3</string>
        <string>/Users/ovc/Projects/LES_v2/sovushka_ng.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ovc/Projects/LES_v2</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TRUSTED_NETWORKS</key>
        <string>127.0.0.1/32,::1/128,10.195.146.0/24</string>
        <key>TRUSTED_NETWORK_ROLE</key>
        <string>admin</string>
        <key>TRUSTED_PROXY_NETWORKS</key>
        <string>127.0.0.1/32,::1/128,10.195.146.136/32</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ovc/Projects/LES_v2/logs/sovushka.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ovc/Projects/LES_v2/logs/sovushka.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
