# sovushka_launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.les.sovushka</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ovc/Projects/LES_v2/sovushka_ng.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ovc/Projects/LES_v2</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ovc/Projects/LES_v2/logs/sovushka.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ovc/Projects/LES_v2/logs/sovushka.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
